# webapp/supervisor-app.conf
[supervisord]
[program:app-uwsgi]
user=uwsgi ;
command = /usr/bin/uwsgi --ini /home/docker/code/uwsgi.ini --mime-file /etc/mime.types ;
